# src/systemd/user/pipetron.service.in
[Unit]
Description=Electron node manager for PipeWire
After=pipewire.service

[Service]
ExecStart=@BINARY_PATH@
Restart=on-failure

[Install]
WantedBy=pipewire.service